# data/ayatana-indicator-session.service.in
[Unit]
Description=Ayatana Indicator Session Service
PartOf=graphical-session.target
After=ayatana-indicators.target

[Service]
ExecStart=@pkglibexecdir@/ayatana-indicator-session-service
Restart=on-failure

[Install]
WantedBy=ayatana-indicators.target
