# data/indicator-session.conf.in
description "Indicator Session Service"

# NOTE: Limiting only to Unity 7 right now as it's only needed
# on the desktop

start on (indicators-loaded or indicator-services-start) and xsession SESSION=ubuntu
stop on desktop-end or indicator-services-end

env G_MESSAGES_DEBUG=all
export G_MESSAGES_DEBUG

respawn

exec @pkglibexecdir@/indicator-session-service
